# cbreader/crontab
*/15    *       *       *       *       python /src/manage.py scan_comics